# ha_opencode/rootfs/etc/s6-overlay/s6-rc.d/ppq-private-proxy/run
#!/command/with-contenv bashio
# ==============================================================================
# OpenCode - Optional PPQ Private Mode Proxy Service
# ==============================================================================

export HOME="/data"
export XDG_DATA_HOME="/data/.local/share"
export XDG_CONFIG_HOME="/data/.config"
export NODE_OPTIONS="--max-old-space-size=256"
export PORT="8787"

if [ -f /data/.env_vars ]; then
    source /data/.env_vars
fi

PPQ_PRIVATE_ENABLED=$(bashio::config 'ppq_private_enabled' 2>/dev/null || echo "false")

if ! bashio::var.true "${PPQ_PRIVATE_ENABLED}"; then
    exec sleep infinity
fi

if [ -z "${PPQ_API_KEY:-}" ]; then
    bashio::log.error "PPQ private mode is enabled, but no PPQ API key is configured."
    bashio::log.error "Set the PPQ API Key option in the add-on Configuration tab, then restart the add-on."
    exec sleep infinity
fi

bashio::log.info "Starting PPQ private-mode proxy on 127.0.0.1:${PORT}"
bashio::log.info "The PPQ proxy is internal-only and is not exposed through Home Assistant networking."

PPQ_PROXY_ENTRYPOINT="$(npm root -g)/ppq-private-mode/bin/server.ts"
exec tsx "${PPQ_PROXY_ENTRYPOINT}"
